# misc/irqbalance.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=/path/to/irqbalance.env
ExecStart=/usr/sbin/irqbalance --foreground $IRQBALANCE_ARGS

[Install]
WantedBy=multi-user.target
