# misc/irqbalance.service
[Unit]
Description=irqbalance daemon
After=syslog.target

[Service]
EnvironmentFile=/etc/sysconfig/irqbalance
ExecStart=/usr/sbin/irqbalance --foreground $ONESHOT

[Install]
WantedBy=multi-user.target
